# test/packaging/amd64/DEBIAN-dnsproxy/postrm
#!/bin/sh

update-rc.d -f dnshipproxy remove 

exit 0



