# test/packaging/amd64/DEBIAN-dnsproxy/prerm
#!/bin/sh

invoke-rc.d --force --quiet dnshipproxy stop 

exit 0


