# test/packaging/debian-init.d-hipd
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:		hipd
# Required-Start:	$remote_fs $syslog
# Required-Stop:	$remote_fs $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		1
# Short-Description:	HIPL IPsec key and mobility management daemon
### END INIT INFO

HIPD_OPTS="-bk"
PID=`pidof -o $$ hipd`

set -e

test -x /usr/sbin/hipd || exit 0

if test -f /etc/default/hipd; then
    . /etc/default/hipd
fi

stop_hipd() {
    # Ensure that we don't try to kill init.d script itself
    if [ $PID ]
    then
	kill $1 $PID
    fi
}

case "$1" in
  start)
        echo -n "Starting hipd"
	start-stop-daemon --start --quiet --exec /usr/sbin/hipd -- $HIPD_OPTS
        echo "."
	;;
  stop)
        echo -n "Stopping hipd"
	stop_hipd
        echo "."
	;;

  reload|force-reload)
        echo -n "Reloading hipd"
	start-stop-daemon --stop --signal 1 --quiet --oknodo --exec /usr/sbin/hipd -name hipd
	echo "."
	;;

  restart)
        echo -n "Restarting hipd"
        # Daemon waits N seconds for closing of connections - just kill it
	stop_hipd -KILL
	start-stop-daemon --start --quiet --exec /usr/sbin/hipd -- $HIPD_OPTS
	echo "."
	;;

  *)
	echo "Usage: /etc/init.d/hipd {start|stop|reload|force-reload|restart}"
	exit 1
esac

exit 0
