# test/packaging/i386/DEBIAN-FW/postinst
#!/bin/sh

update-rc.d hipfw defaults 
echo "invoke-rc.d --quiet hipfw start" | at now + 1 min 2>/dev/null 
echo "hipd firewall will start"

