# test/packaging/i386/DEBIAN-dnsproxy/postrm
#!/bin/sh

update-rc.d -f dnshipproxy remove 



