# test/packaging/i386/DEBIAN-dnsproxy/preinst
#!/bin/sh

invoke-rc.d --force --quiet dnshipproxy stop 
update-rc.d -f dnshipproxy remove 

exit 0 

 
