# test/packaging/i386/DEBIAN/postrm
#!/bin/sh
update-rc.d -f hipfw remove 
update-rc.d -f hipd remove 


