# test/packaging/i386/DEBIAN/prerm
#!/bin/sh
invoke-rc.d --force --quiet hipfw stop 
invoke-rc.d --force --quiet hipd stop 


