# test/packaging/amd64/DEBIAN/postinst
#!/bin/sh

#ldconfig

invoke-rc.d --force --quiet hipd stop
update-rc.d -f hipd remove
update-rc.d hipd defaults
invoke-rc.d hipd start
